# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:webapp]
directory=/app
command=/app/webapp
